# .config/systemd/user/niri-session.service
[Unit]
Description=Start Niri Session

[Service]
ExecStart=/usr/bin/niri-session
Restart=always
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
